# root/etc/services.d/spyserver/run
#!/usr/bin/with-contenv bash

exec /app/spyserver /config/spyserver.config
