# root/etc/s6-overlay/s6-rc.d/spyserver/run
#!/command/with-contenv bash

if [ -n "$SKIP_APP" ]; then
	exit 0
fi

exec /app/spyserver /config/spyserver.config
